# scrape-and-tweet/crontab
* * * * * /usr/local/bin/node /app/scrape-and-tweet/index.js > /dev/stdout
